# installer/linux/torvm.service
[Unit]
Description=TorVM - Transparent Tor Proxy Virtual Machine
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/torvm --headless
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=10
WatchdogSec=60
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
